# reactphp/resources/supervisord.conf
[supervisord]
nodaemon=true

[program:react]
command=/usr/bin/php /var/www/html/public/index.php

[program:php-fpm]
command=/usr/sbin/php-fpm7.4 -F
